# AirplanesLiveHA/rootfs/etc/services.d/airplanes_live_api/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the Airplanes Live API services (Web Interface + MQTT)
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

bashio::log.info "Starting Airplanes Live API services"

# Start the web interface in the background
bashio::log.info "Starting web interface..."
python3 /app/app.py &

# Wait a moment for the web interface to start
sleep 2

# Start the MQTT service in the foreground
bashio::log.info "Starting MQTT service..."
exec python3 /app/run.py 